# supervisor-queue.conf
[program:miysoft-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/queue-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:miysoft-queue-exports]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --queue=exports --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/queue-exports.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:miysoft-queue-reports]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --queue=reports --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/queue-reports.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:miysoft-queue-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --queue=notifications --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/queue-notifications.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:miysoft-queue-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --queue=emails --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/queue-emails.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:miysoft-scheduler]
process_name=%(program_name)s
command=php /var/www/miysoft_pts/artisan schedule:run
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/miysoft_pts/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5