# docs/supervisor_queue.conf
; MİYSOFT PTS — Supervisor Queue Worker Configuration
; Kullanım: sudo ln -s /var/www/miysoft_pts/docs/supervisor_queue.conf /etc/supervisor/conf.d/miysoft_queue.conf
; Sonra: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start miysoft-queue:*

[program:miysoft-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/miysoft_pts/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default,notifications,exports
directory=/var/www/miysoft_pts
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
startretries=3
stderr_logfile=/var/log/miysoft_queue_error.log
stdout_logfile=/var/log/miysoft_queue.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:miysoft-queue-schedule]
command=php /var/www/miysoft_pts/artisan schedule:work
directory=/var/www/miysoft_pts
autostart=true
autorestart=true
user=www-data
stderr_logfile=/var/log/miysoft_schedule_error.log
stdout_logfile=/var/log/miysoft_schedule.log
stdout_logfile_maxbytes=5MB
